# config/mvt.cron
0 * * * * /data/scripts/pullUpdatedSTVs.sh
